# config/systemd/user/home-tonya-Music.automount
[Unit]
  Description=cifs mount script
  Requires=network-online.target
  After=network-online.service

[Automount]
  Where=/home/tonya/Music

[Install]
  WantedBy=multi-user.target
