# config/systemd/user/home-tonyampomah-mnt-Downloads.automount
[Unit]
  Description=Automount Downloads

[Automount]
  Where=/home/tonyampomah/mnt/Downloads

[Install]
  WantedBy=multi-user.target
