# config/systemd/user/home-kwame-dat-Pictures.automount
[Unit]
  Description=Automount Pictures

[Automount]
  Where=/home/kwame-dat/Pictures

[Install]
  WantedBy=multi-user.target
