# config/systemd/user/home-tonya-mnt-Downloads.automount
[Unit]
  Description=Automount Downloads

[Automount]
  Where=/home/tonya/mnt/Downloads

[Install]
  WantedBy=multi-user.target
