# config/systemd/user/home-tonya-Pictures.automount
[Unit]
  Description=Automount Pictures

[Automount]
  Where=/home/tonya/Pictures

[Install]
  WantedBy=multi-user.target
