# config/systemd/user/home-kwame-dat-mnt-Downloads.automount
[Unit]
  Description=Automount Downloads

[Automount]
  Where=/home/kwame-dat/mnt/Downloads

[Install]
  WantedBy=multi-user.target
