# config/systemd/user/home-tonya-mnt-Iso.automount
[Unit]
  Description=Automount Iso

[Automount]
  Where=/home/tonya/mnt/Iso

[Install]
  WantedBy=multi-user.target
