# config/systemd/user/home-kwame-dat-mnt-Iso.automount
[Unit]
  Description=Automount Iso

[Automount]
  Where=/home/kwame-dat/mnt/Iso

[Install]
  WantedBy=multi-user.target
