# config/systemd/user/home-kwame-dat-Videos.automount
[Unit]
  Description=Automount Videos

[Automount]
  Where=/home/kwame-dat/Videos

[Install]
  WantedBy=multi-user.target
