# config/systemd/user/home-tonyampomah-Videos.automount
[Unit]
  Description=Automount Videos

[Automount]
  Where=/home/tonyampomah/Videos

[Install]
  WantedBy=multi-user.target
