# config/systemd/user/home-tonya-Videos.automount
[Unit]
  Description=Automount Videos

[Automount]
  Where=/home/tonya/Videos

[Install]
  WantedBy=multi-user.target
