# config/systemd/user/home-tonyampomah-mnt-Iso.automount
[Unit]
  Description=Automount Iso

[Automount]
  Where=/home/tonyampomah/mnt/Iso

[Install]
  WantedBy=multi-user.target
