# config/systemd/user/home-tonyampomah-Pictures.automount
[Unit]
  Description=Automount Pictures

[Automount]
  Where=/home/tonyampomah/Pictures

[Install]
  WantedBy=multi-user.target
