# config/systemd/user/mnt-Videos.automount
[Unit]
  Description=Automount Videos

[Automount]
  Where=/mnt/Videos

[Install]
  WantedBy=multi-user.target
