# config/systemd/user/mnt-Pictures.automount
[Unit]
  Description=Automount Pictures

[Automount]
  Where=/mnt/Pictures

[Install]
  WantedBy=multi-user.target
