# config/systemd/user/mnt-Iso.automount
[Unit]
  Description=Automount Iso

[Automount]
  Where=/mnt/Iso

[Install]
  WantedBy=multi-user.target
