# config/systemd/user/mnt-Downloads.automount
[Unit]
  Description=Automount Downloads

[Automount]
  Where=/mnt/Downloads

[Install]
  WantedBy=multi-user.target
