# .config/systemd/home-tonya-Videos.automount
[Unit]
  Description=cifs mount script
  Requires=network-online.target
  After=network-online.service

[Automount]
  What=//192.168.1.3/Videos
  Where=/home/tonya/Videos
  Options=_netdev,credentials=/etc/samba/credentials/music-share,iocharset=utf8,vers=2.0,file_mode=0777
  Type=cifs

[Install]
  WantedBy=multi-user.target
